# root/etc/services.d/websocketd-script2-task/run
#!/command/with-contenv sh

echo "starting script2 task websocketd"
exec /usr/bin/websocketd -port=8085 -address=127.0.0.1 /opt/scripts/script2-task

