# root/etc/services.d/websocketd-script4/run
#!/command/with-contenv sh

echo "starting script4 websocketd"
exec /usr/bin/websocketd -port=8083 -address=127.0.0.1 /opt/scripts/script4

