# root/etc/periodic/15min/azure-logs
#!/bin/bash

grep -i azure /var/log/livecd.log > /var/log/azure.log


