# root/etc/services.d/websocketd-script4-task/run
#!/command/with-contenv sh

echo "starting script4-task websocketd"
exec /usr/bin/websocketd -port=8087 -address=127.0.0.1 /opt/scripts/script4-task

